# packaging/wardragon-console.service
[Unit]
Description=WarDragon Console
After=network.target zmq-decoder.service wardragon-monitor.service dragonsync.service
Wants=zmq-decoder.service wardragon-monitor.service dragonsync.service

[Service]
Type=simple
User=dragon
Group=dragon
Environment=WARDRAGON_CONSOLE_HOST=127.0.0.1
Environment=WARDRAGON_CONSOLE_PORT=4280
Environment=WARDRAGON_CONSOLE_CONFIG_WRITE=1
Environment=WARDRAGON_CONSOLE_REMOTE_CONFIG_WRITE=1
Environment=WARDRAGON_CONSOLE_REMOTE_RESTART=1
Environment=WARDRAGON_CONSOLE_TETHER_ENABLED=1
Environment=WARDRAGON_CONSOLE_TETHER_PORT=4280
Environment=WARDRAGON_CONSOLE_TETHER_CIDRS=192.168.42.0/24,192.168.43.0/24,172.20.10.0/28
Environment=WARDRAGON_DRAGONSCOPE_DIR=/home/dragon/WarDragon/dragonsdr_dji_droneid
ExecStart=/usr/local/bin/wardragon-console
Restart=on-failure
RestartSec=5

# Hardening. The directives below are deliberately limited to ones that do
# NOT implicitly enable NoNewPrivileges. The console uses sudo to restart
# dragonsync.service and to add the stable tether IP alias, both through
# narrow sudoers rules; sudo refuses to run under NoNewPrivileges=yes, and
# many systemd hardening directives (ProtectKernelTunables=, RestrictNamespaces=,
# LockPersonality=, etc.) imply that flag whether or not we set it explicitly.
# The trust model here is "physical/local-access = trust"; the kept directives
# still meaningfully constrain filesystem and namespace access.
NoNewPrivileges=no
PrivateTmp=true

[Install]
WantedBy=multi-user.target
